# misc/etc/cron.d/colab
*   *   *   *   *   colab   colab-admin update_index --age=1 &> /dev/null

34  2   *   *   *   colab   colab-admin rebuild_index --noinput &> /dev/null

*   1   *   *   *   colab   colab-admin cleanup_snippets &> /dev/null

*/5 *   *   *   *   colab   colab-admin update_badges &> /dev/null

27  3   *   *   *   colab   colab-admin clearsessions &> /dev/null

10  *   *   *   *   colab   colab-admin update_planet &> /dev/null

*   *   *   *   *   colab   colab-admin import_emails --archives_path=/var/lib/mailman/archives/private/ &> /dev/null
