# misc/etc/cron.d/colab
*   *   *   *   *   colab   colab-admin update_index --age=1 &> /dev/null

34  2   *   *   *   colab   colab-admin rebuild_index --noinput &> /dev/null

27  3   *   *   *   colab   colab-admin clearsessions &> /dev/null

*   *   *   *   *   colab   colab-admin import_emails &> /dev/null
