# misc/lib/systemd/system/colab.service
[Unit]
Description=Colab
Before=celeryd.service celerybeat.service
Requires=celeryd.service celerybeat.service

[Service]
User=colab
EnvironmentFile=-/etc/sysconfig/colab
PermissionsStartOnly=true
ExecStartPre=/usr/bin/install -d -m 0750 -o colab -g colab /var/lock/colab
ExecStart=/usr/lib/colab/bin/gunicorn colab.wsgi:application --config=/etc/colab/gunicorn.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
