# puppet/modules/supervisor/templates/supervisor.conf.erb
[program:<%= @app_name %>]
command=<%= @command %>
directory=<%= @directory %>
user=<%= @user %>
autostart=true
autorestart=true
redirect_stderr=True
<%- if @startsecs -%>
startsecs=<%= @startsecs %>
<%- end -%>
<%- if @stopwaitsecs -%>
stopwaitsecs=<%= @stopwaitsecs %>
<%- end -%>
<%- if @priority -%>
priority=<%= @priority %>
<%- end -%>
